# deployment/meditrack.service
# ===========================================================
# systemd service for MediTrack (Gunicorn)
#
# Install on the EC2 instance:
#   sudo cp deployment/meditrack.service /etc/systemd/system/meditrack.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable meditrack
#   sudo systemctl start meditrack
#   sudo systemctl status meditrack
#
# Adjust User, WorkingDirectory, and the venv path below to match
# your actual EC2 username and clone location.
# ===========================================================

[Unit]
Description=Gunicorn instance serving the MediTrack Flask app
After=network.target

[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/capstone_sdp
EnvironmentFile=/home/ubuntu/capstone_sdp/.env
ExecStart=/home/ubuntu/capstone_sdp/venv/bin/gunicorn \
    --workers 3 \
    --bind 127.0.0.1:5000 \
    app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
